# systemd.sample/trojan.service
[Unit]
Description=Trojan net

[Service]
Type=simple
ExecStart=trojan -c /etc/trojan/client.config.json

[Install]
WantedBy=multi-user.target